# services/xai_rt_proxy/sola-xai-rt-proxy.service
[Unit]
Description=SOLA xAI Realtime WebSocket proxy
After=network.target

[Service]
Type=simple
User=sola
Group=sola
WorkingDirectory=/opt/sola-xai-rt-proxy
EnvironmentFile=/etc/sola-xai-rt-proxy.env
ExecStart=/usr/bin/node /opt/sola-xai-rt-proxy/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sola-xai-rt-proxy

# Sandboxing: the proxy needs only outbound TLS and a local bind port.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
